# media-center.service
[Unit]
Description=Media Center Dashboard
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStartPre=/bin/sleep 5
ExecStart=/home/user/.local/bin/media-center.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
